# config/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.

* * * * * /usr/bin/php /var/www/artisan schedule:run > /dev/null 2>&1
